# deployment/mdai-system.service
[Unit]
Description=MDai Liveness Detection System
Documentation=https://github.com/yourusername/mdai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/mercleDev/codebase
ExecStart=/home/mercleDev/codebase/build/mdai_system
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Peripheral requirements
Environment="DISPLAY=:0"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Wait for USB devices
ExecStartPre=/bin/sleep 5

# Ensure permissions for all ttyACM devices
ExecStartPre=/bin/bash -c 'for dev in /dev/ttyACM*; do [ -e "$dev" ] && chmod 666 "$dev"; done || true'

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Watchdog (restart if unresponsive)
WatchdogSec=60

# Kill everything on stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
