# deployment/mdai-system.service
[Unit]
Description=MDai Liveness Detection System
After=local-fs.target
After=NetworkManager.service
After=wifi-always-on.service
After=tailscaled.service
After=ssh.service

Wants=wifi-always-on.service
Wants=tailscaled.service

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=forking
User=mercleDev
Group=mercleDev

# Start EXACTLY like rebuild-and-restart.sh does
ExecStartPre=-/usr/bin/pkill -9 -x mdai_system
ExecStartPre=/bin/sleep 2

ExecStart=/bin/bash -c 'cd /home/mercleDev/codebase && nohup ./build/mdai_system >> /home/mercleDev/mdai_logs/rdk.log 2>&1 &'

# Give it time to fork
ExecStartPost=/bin/sleep 3

Restart=always
RestartSec=5

KillMode=process
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
