# deployment/tailscale-dns-monitor.service
[Unit]
Description=Tailscale DNS Health Monitor
After=tailscaled.service systemd-resolved.service
Wants=tailscaled.service

[Service]
Type=simple
ExecStart=/usr/local/bin/tailscale-dns-monitor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
