# rdk-service/mdai-system.service
[Unit]
Description=MDAI System - Primary RDK Service
Documentation=file:///home/mercleDev/codebase/README.md

# ============================================
# BOOT ORDER - After network is ready
# ============================================
After=network-online.target
After=NetworkManager.service
After=wpa_supplicant.service
After=tailscaled.service
After=systemd-udevd.service
Wants=network-online.target

# No limit on restart attempts
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=simple

# ============================================
# RUN AS ROOT - Required for real-time scheduling
# ============================================
User=root
Group=root

# Hardware access (camera, serial, USB)
SupplementaryGroups=dialout video plugdev input

# ============================================
# ENVIRONMENT & PATHS
# ============================================
WorkingDirectory=/home/mercleDev/codebase/build
Environment="LD_LIBRARY_PATH=/home/mercleDev/codebase/build:/home/mercleDev/codebase/lib/mediapipe:/usr/local/lib"
Environment="HOME=/home/mercleDev"
Environment="XDG_RUNTIME_DIR=/run/user/0"

# ============================================
# STARTUP
# ============================================
# Create log directory
ExecStartPre=/bin/mkdir -p /home/mercleDev/mdai_logs

# SINGLE INSTANCE ENFORCEMENT: Kill any orphaned/manual instances
ExecStartPre=/bin/bash -c 'pkill -9 -f "mdai_system$" 2>/dev/null || true'

# Small delay to ensure hardware is ready and old process is gone
ExecStartPre=/bin/sleep 2

# Main executable
ExecStart=/home/mercleDev/codebase/build/mdai_system

# PID file for tracking
PIDFile=/run/mdai-system.pid

# ============================================
# CPU PRIORITY - MAXIMUM
# ============================================
# Highest nice value (-20 is max priority)
Nice=-20

# Best-effort I/O with highest priority
IOSchedulingClass=best-effort
IOSchedulingPriority=0

# No CPU throttling
# CPUQuota removed - 'infinity' is not valid

# Use all 8 cores
CPUAffinity=0-7

# Maximum CPU weight (default is 100)
CPUWeight=10000

# ============================================
# MEMORY - MAXIMUM PRIORITY
# ============================================
# NEVER kill this process for OOM
OOMScoreAdjust=-1000
OOMPolicy=continue

# Reserve minimum 512MB, prefer up to 4GB
MemoryMin=512M
MemoryHigh=4G

# Lock memory to prevent swapping (no swap anyway)
LockPersonality=true

# ============================================
# RESTART POLICY - NEVER STAY DOWN
# ============================================
Restart=always
RestartSec=2

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=30

# ============================================
# SHUTDOWN
# ============================================
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
SendSIGKILL=yes

# ============================================
# WATCHDOG - Disabled (app doesn't support sd_notify)
# ============================================
# WatchdogSec disabled - app doesn't send heartbeats
# We rely on Restart=always + external mdai-watchdog service instead

# ============================================
# LOGGING
# ============================================
StandardOutput=append:/home/mercleDev/mdai_logs/system.log
StandardError=append:/home/mercleDev/mdai_logs/system.log
SyslogIdentifier=mdai-system

# ============================================
# SECURITY (minimal - need hardware access)
# ============================================
PrivateTmp=false
PrivateNetwork=false
ProtectSystem=false
ProtectHome=false

# Allow all capabilities for hardware access
AmbientCapabilities=CAP_SYS_NICE CAP_SYS_RAWIO CAP_NET_RAW CAP_IPC_LOCK
CapabilityBoundingSet=~

[Install]
WantedBy=multi-user.target
